# services/ft-nowplaying/ft-nowplaying.service
[Unit]
Description=KNOB now-playing display on FlaschenTaschen
After=network.target liquidsoap.service
Wants=liquidsoap.service

[Service]
Type=simple
User=nthmost
WorkingDirectory=/home/nthmost
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/nthmost/ft_nowplaying.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
